# rootfs/etc/services.d/99-flood/run
#!/usr/bin/with-contenv sh

export $FLOOD_AUTH
s6-setuidgid abc flood \
    -h localhost -p 3000 -d /config/flood \
    --rtsocket /run/php/.rtorrent.sock \
    --auth ${FLOOD_AUTH:-none}